# recipes-tweaks/dnsmasq/dnsmasq/99-dnsmasq.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="br-lan", RUN+="/etc/init.d/dnsmasq start"
